# etc/insserv/override/mountall.sh
#! /bin/sh
### BEGIN INIT INFO
# Provides:          mountall
# Required-Start:    checkfs checkroot-bootclean
# Required-Stop: 
# Default-Start:     
# Default-Stop:
# Short-Description: Mount all filesystems.
# Description:
### END INIT INFO
