# etc/insserv/overrides/hwclock.sh
#!/bin/sh
### BEGIN INIT INFO
# Provides:          hwclock
# Required-Start:    
# Required-Stop:     $local_fs
# Default-Start:     
# X-Start-Before:    checkroot
# Default-Stop:      0 6
### END INIT INFO
