# etc/insserv/overrides/mtab.sh
#! /bin/sh
### BEGIN INIT INFO
# Provides:          mtab
# Required-Start:    init
# Required-Stop:
# Default-Start:     S
# Default-Stop:
# Short-Description: Update mtab file.
# Description:       Update the mount program's mtab file after
#                    all local filesystems have been mounted.
### END INIT INFO
