# etc/init.d/sleep.sh
#! /bin/bash
### BEGIN INIT INFO
# Provides:	      	         sleep
# Required-Start:	
# Required-Stop:
# Should-Start:      
# Default-Start:     
# Default-Stop:
# Short-Description: 
# Description:  Sleep for a second
### END INIT INFO
sleep 1
